# root/etc/systemd/system/dualgpio.service
[Unit]
Description=GPIO Button Watcher for Dual MPV
After=dualmp4.service
Requires=dualmp4.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/gpio_switch.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target