# root/etc/systemd/system/mpvsync.service
[Unit]
Description=MPV dual sync watchdog
After=dualmp4.service
Requires=dualmp4.service
PartOf=dualmp4.service

[Service]
Type=simple
User=pi
Group=pi
ExecStart=/usr/bin/python3 /boot/firmware/splash/mpv_sync_watchdog.py
Restart=always
RestartSec=3
TimeoutStopSec=5

[Install]
WantedBy=dualmp4.service
