# root/etc/systemd/system/mpvsync.service
[Unit]
Description=MPV Sync Watchdog
After=dualmp4.service
Requires=dualmp4.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/mpv_sync_watchdog.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target